# network/if-pre-up.d/nat-redirect-80-443
#!/bin/sh

iptables-restore < /etc/network/redirect-80-443.iptables
